# root/etc/s6-overlay/s6-rc.d/init-mod-route2me-install-deps/up
/etc/s6-overlay/s6-rc.d/init-mod-route2me-install-deps/run